# Pixie/audio/systemd/ensure-master-volume.timer
[Unit]
Description=Run ensure-master-volume.service every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target 