# scripts/pi_access_point/setup_wifi_ap.service
[Unit]
Description=Setup WiFi AP (ap0) on boot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup_wifi_ap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target